# Dockerfile/loopbackapi/configs/services.ini
#[program:sshd]
#command=/usr/sbin/sshd -D

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:php-fpm]
command=/usr/sbin/php-fpm --daemonize
